# xqlist.cron
###fetch xueqie stock list info
15 9 * * 2,3,4,5,6 root /home/ray/workspace/uss-practice/fetchXueqiuList.py >>/tmp/fxq.log 2>>/tmp/fxq.err

